# .config/systemd/user/lighttpd.service
[Unit]
Description=Lighttpd Web Server
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/lighttpd-angel -D -f /home/ayekat/.config/lighttpd/lighttpd.conf
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
