# etc/karuiwm/profile.d/80-nitrogen.sh
if command -v nitrogen >/dev/null; then
	nitrogen --restore
fi
