# etc/karuiwm/profile.d/90-scratchpad.sh
if have urxvt tmux; then
	urxvt -e "$XDG_LIB_HOME"/tmux/tmux-view SCRATCHPAD &
fi
