# etc/karuiwm/profile.d/95-browser.sh
# Prevent xdg-open from launching Firefox everytime something fails with
# qutebrowser:
export BROWSER=qutebrowser
