# .local/etc/systemd/user/lighttpd.service
[Unit]
Description=Lighttpd Web Server
After=syslog.target network.target

[Service]
ExecStartPre=/usr/bin/lighttpd -t -f ${XDG_CONFIG_HOME}/lighttpd/lighttpd.conf
ExecStart=/usr/bin/lighttpd-angel -D -f ${XDG_CONFIG_HOME}/lighttpd/lighttpd.conf
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=default.target
