# etc/systemd/user/error.service
[Unit]
Description=Test service that fails

[Service]
ExecStart=/bin/false
