# etc/systemd/user/mbsync@.service
[Unit]
Description=Mailbox synchronisation service for %I

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -c ${XDG_CONFIG_HOME}/mbsync/mbsyncrc %i
