# etc/systemd/user/talky-ssh-agent.service
[Unit]
Description=SSH agent with notifications

[Service]
ExecStart=%h/.local/bin/talky-ssh-agent

[Install]
WantedBy=default.target
