# etc/sh/profile.d/30-utils.sh
# utils

path prepend "$XDG_LIB_HOME/utils/bin"
