# etc/sh/profile.d/40-salt.sh
# SaltStack

path prepend "$XDG_LIB_HOME"/salt/bin
