# etc/sh/profile.d/60-editor.sh
# Text editors/pagers

export EDITOR='/usr/bin/vim'
export VISUAL='/usr/bin/vim'
export PAGER='/usr/bin/less'
