# etc/sh/profile.d/dotfiles.sh
# dotfiles

if [ -d "$HOME/.local/lib/dotfiles/bin" ]; then
	export PATH="$HOME/.local/lib/dotfiles/bin:$PATH"
fi
