# etc/sh/profile.d/dotfiles.sh
# dotfiles

path prepend "$XDG_LIB_HOME/dotfiles/wrappers"
