# etc/sh/profile.d/dotfiles.sh
# dotfiles

if [ -d "$XDG_LIB_HOME/dotfiles/bin" ]; then
	export PATH="$XDG_LIB_HOME/dotfiles/bin:$PATH"
fi
