# etc/sh/profile.d/utils.sh
# utils

if [ -d "$XDG_LIB_HOME/utils/bin" ]; then
	export PATH="$XDG_LIB_HOME/utils/bin:$PATH"
fi
