# .local/etc/sh/profile.d/dotfiles.sh
# dotfiles

if [ -d "$HOME/.local/bin/dotfiles" ]; then
	export PATH="$HOME/.local/bin/dotfiles:$PATH"
fi
