# .local/etc/sh/profile.d/epflscripts.sh
# EPFL scripts

if [ -d "$HOME/.local/bin/epflscripts" ]; then
	export PATH="$PATH:$HOME/.local/bin/epflscripts"
fi
