# .local/etc/sh/profile.d/utils.sh
# utils

if [ -d "$HOME/.local/bin/utils" ]; then
	export PATH="$HOME/.local/bin/utils:$PATH"
fi
