# .local/etc/sh/profile.d/utils.sh
# utils

if [ -d "$HOME/.local/lib/utils/bin" ]; then
	export PATH="$HOME/.local/lib/utils/bin:$PATH"
fi
